# crontab
@yearly /home/corvusoft/annual-process

*/10 * * * * /home/corvusoft/fschk